# turn_tv_off.timer
[Unit]
Description=Turn TV off at specific times

[Timer]
OnCalendar=Mon..Thu 17:30
OnCalendar=Fri 16:30

[Install]
WantedBy=timers.target
