# services/start_vscode_tunnel.service
[Unit]
Description=VS Code Tunnel Service
After=network.target

[Service]
ExecStart=code tunnel
User=pauls

[Install]
WantedBy=multi-user.target