# services/turn_tv_on.timer
[Unit]
Description=Turn TV on at specific times

[Timer]
OnCalendar=Mon..Wed 5:30

[Install]
WantedBy=timers.target
